# build/etc/systemd/system/resource_monitor.service
[Unit]
Description=Resource Monitoring Daemon
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/bin/resource-monitor
Restart=always
User=ivan
Group=ivan

[Install]
WantedBy=multi-user.target
